# conf/etc/services.d/01-synapse/run
#!/bin/sh
NAME="SYNAPSE" 
COMMAND="/usr/local/bin/synapse-launcher" 
WAIT="/tmp/synapse"
exec service-handler "$NAME" "$COMMAND" "$WAIT"